# build-armbian/armbian-files/different-files/lpa3399pro/rootfs/etc/udev/rules.d/99-lpa-hdmi-unblank.rules
ACTION=="add", SUBSYSTEM=="graphics", KERNEL=="fb0", RUN+="/bin/sh -c 'echo 0 > /sys/class/graphics/fb0/blank'"
ACTION=="change", SUBSYSTEM=="drm", KERNEL=="card0-*", RUN+="/bin/sh -c 'echo 0 > /sys/class/graphics/fb0/blank 2>/dev/null || true'"
