# build-armbian/armbian-files/different-files/lpa3399pro/rootfs/etc/systemd/system/docker.service.d/99-tf-data-root.conf
[Unit]
# Optional TF data disk at /mnt/sdcard (fstab nofail). Without TF, dockerd does not start.
After=mnt-sdcard.mount
Wants=mnt-sdcard.mount
ConditionPathIsMountPoint=/mnt/sdcard
