# build-armbian/armbian-files/different-files/lpa3399pro/rootfs/etc/systemd/system/lpa-hdmi-unblank.service
[Unit]
Description=LPA3399Pro HDMI framebuffer unblank + tty1 kick
After=local-fs.target systemd-modules-load.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/lpa-hdmi-unblank.sh
RemainAfterExit=yes
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
